# _playground/_dev/packages/bme680-service/services/bme680-heat-soak.service
[Unit]
Description=BME680 Heat Soak Detection Service
After=network.target mqtt.service
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=%h/.local/share/bme680-service
Environment="PATH=%h/.local/share/bme680-service/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=%h/.local/share/bme680-service/.venv/bin/python %h/.local/share/bme680-service/monitor-temperature.py --monitor --mqtt
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bme680-heat-soak

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/share/bme680-service

[Install]
WantedBy=multi-user.target
