# _playground/_dev/packages/bme680-service/services/bme680-heatsoak.service
[Unit]
Description=BME680 HeatSoak Consumer
After=network-online.target mosquitto.service bme680-base.service
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/.local/share/bme680-service
Environment="PATH=/home/pi/.local/share/bme680-service/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/pi/.local/share/bme680-service/.venv/bin/python /home/pi/.local/share/bme680-service/heatsoak-consumer.py --topic-raw sensors/bme680/raw --topic-state sensors/bme680/heatsoak/state --target 45.0 --plateau-slope 0.05 --window 120
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/.local/share/bme680-service

[Install]
WantedBy=multi-user.target


